# twentynine/otus.service
[Unit]

Description=Runs /usr/local/bin/mystartup.sh

  
[Service]

ExecStart=/etc/rc.local 


[Install]

WantedBy=multi-user.target

